# services.d/cockroach/run
#!/bin/bash

cockroach start \
  --insecure \
  --store=$(hostname) \
  --listen-addr=0.0.0.0:26257 \
  --http-addr=0.0.0.0:8080 \
  --join=localhost:26257